# patches/etc/init.d/oc-ssh-keys
#!/sbin/runscript

description="Sync SSH keys"


start() {
    ebegin "Syncing SSH keys"
    mkdir -p /root/.ssh/
    chmod 700 /root/.ssh
    /usr/local/bin/oc-metadata --cached | grep SSH_PUBLIC_KEYS_.*_KEY | cut -d'=' -f 2- | sed "s/^'//;s/'$//" | tr -d \\ > /root/.ssh/authorized_keys
    eend $?
}
